# systemd/imonvpn-awg-traffic.timer
[Unit]
Description=Collect imonvpn AmneziaWG traffic counters every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
